# etc/systemd/system/cloudflare-ddns.timer
[Unit]
Description=cloudflare-ddns

[Timer]
OnBootSec=30
OnUnitActiveSec=300

[Install]
WantedBy=timers.target
