# etc/systemd/system/cloudflare-ddns.service
[Unit]
Description=cloudflare-ddns

[Service]
EnvironmentFile=/root/.cloudflare-ddns
ExecStart=/usr/local/bin/cloudflare-ddns
